# supervisor/supervisord.d/mongodb.ini
[program:mongodb]
command=/usr/local/mongodb/bin/mongod --config /usr/local/mongodb/mongodb.conf --directoryperdb --quiet --profile 1 --slowms 1000 --noprealloc --shardsvr
directory=/data/mongodb
autostart=true
user=root
