# supervisor/supervisord.d/mysql.ini
[program:mysql]
command=/usr/sbin/mysqld --pid-file=/var/run/mysqld/mysqld.pid
user=mysql
autostart=true
exitcodes=0,2
