# supervisor/supervisord.d/nginx.ini
[program:nginx]
command=/usr/local/nginx/sbin/nginx
stopsignal=QUIT

